# dockerfiles/php/config/supervisor/php.conf
[program:php]
command=php-fpm
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
